# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true

[program:java]
command=java -jar /usr/src/beru/beru.jar  # Если нужен другой порт, укажите через параметры приложения
autostart=true
autorestart=true